# service/wm8960-audio.service
[Unit]
Description=WM8960 Audio HAT PLL Configuration
After=sound.target
Requires=sound.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wm8960-pll-config.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
